# etc/supervisorx.conf
[supervisord]
nodaemon=true

[program:tmal]
command=/usr/bin/tmax run -config=/etc/tmal/cfx
autorestart=true
priority=200

[supervisorctl]
serverurl = unix:///tmp/supervisor.sock
